# autostart/chrome.desktop
[Desktop Entry]
Type=Application
Exec=chromium
